# etc/services.d/webpack-dev-server/run
#!/usr/bin/with-contenv sh

cd /code
exec /node_modules/.bin/webpack-dev-server
