# etc/services.d/cron/run
#!/usr/bin/with-contenv sh

cd /code
exec mcron mcron.guile
