# etc/services.d/mimir-dev-server/run
#!/usr/bin/with-contenv sh

cd /code
exec mserve --app-name development --server-name development --logging-config development config.toml.mako
